# infra/systemd/blind-nav.service
[Unit]
Description=Blind Navigation System
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=blindnav
ExecStart=/usr/bin/python3 /opt/blind-nav/main.py /etc/blind-nav/config.yaml
Restart=always
RestartSec=5
WatchdogSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=blind-nav

[Install]
WantedBy=multi-user.target
